# launchd/com.pea-tracker.market-watcher-18h.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pea-tracker.market-watcher-18h</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yousrids/Documents/pea-tracker/scripts/run-market-watcher.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>1</integer> <!-- Lundi -->
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>2</integer> <!-- Mardi -->
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>3</integer> <!-- Mercredi -->
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>4</integer> <!-- Jeudi -->
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
            <key>Weekday</key>
            <integer>5</integer> <!-- Vendredi -->
        </dict>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yousrids/Documents/pea-tracker</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/yousrids</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/yousrids/Documents/pea-tracker/logs/launchd-18h.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yousrids/Documents/pea-tracker/logs/launchd-18h.err</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
